# services/ovos-ssh-signal.service
[Service]
ExecStartPre=/usr/bin/python /usr/libexec/ovos-ssh-enabled-signal
ExecStopPost=/usr/bin/python /usr/libexec/ovos-ssh-disabled-signal
