# services/ovos-clock-sync.service
[Service]
ExecStartPost=%h/.venvs/ovos/bin/python /usr/libexec/ovos-clock-sync
