# services/ovos-clock-sync.service
[Service]
ExecStartPost=/bin/bash -c 'nohup /usr/bin/python /usr/libexec/ovos-clock-sync > /dev/null 2>&1 &'