# services/ovos-restart-signal.service
[Service]
ExecStartPre=%h/.venvs/ovos/bin/python /usr/libexec/ovos-restart-signal
