# services/ovos-clock-sync.service
[Service]
ExecStartPost=/usr/bin/python /usr/libexec/ovos-clock-sync
