# overlays/base_ovos/etc/systemd/system/ovos-shutdown-signal.service
[Unit]
Description=Emit OVOS shutdown signal
DefaultDependencies=no
Before=halt.target shutdown.target

[Service]
Type=oneshot
ExecStart=/usr/bin/python /usr/libexec/ovos-shutdown-signal

[Install]
WantedBy=halt.target shutdown.target
