# tuning/combine_sinks.service
[Unit]
Description=Sound output configuration service
PartOf=ovos.service
After=ovos-audio.service

[Service]
Type=oneshot
Group=ovos
UMask=002
ExecStart=/usr/libexec/update-audio-sinks.sh
RemainAfterExit=yes
Restart=on-failure
RestartSec=5

[Install]
WantedBy=ovos.service